# .config/systemd/user/backup-personal.service
[Unit]
Description=Backup Documents to External Drive

[Service]
Type=oneshot
ExecStart=/home/anuradhe/.local/bin/backup_personal.sh

